# vcontrold.service
[Unit]
Description=vcontrold
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/vcontrold

[Install]
WantedBy=multi-user.target
